# hashlatch-stratum.service
[Unit]
Description=HashLatch Services (stratum, api, admin, explorer, tunnel)
After=network.target hashlatch.service
Wants=hashlatch.service
StartLimitIntervalSec=0

[Service]
Type=oneshot
RemainAfterExit=yes
User=dstrychalski
ExecStart=/home/dstrychalski/hashlatch-services.sh
ExecStop=/bin/bash -c 'pkill -f "node server.js"; pkill -f hlc_explorer; pkill -f hashlatch_api; pkill -f admin_panel; pkill -f cloudflared'
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
